# data/crontab
*/10 * * * *    root       bash /root/uniearn/data/cron-fetch.sh
